# nfc-tools/supervisord.conf
[program:pcscd]
command=pcscd --foreground
